# etc/tomcat/bin/tomee
#!/bin/bash
# description: TomEE Start Stop Restart
# processname: tomee
# chkconfig: 234 20 80

TOMEE_USER=iot.admin

if [ ! -v INSTALLATION_HOME ] #INSTALLATION_HOME not set
then
    INSTALLATION_HOME=/opt/tomcat
    export INSTALLATION_HOME
fi

if [ ! -v JAVA_HOME ] #JAVA_HOME not set
then
    echo "JAVA_HOME is not set, set JAVA_HOME and try again."
fi

CATALINA_HOME=$INSTALLATION_HOME/apache-tomee-jaxrs-1.7.4
export CATALINA_HOME

case $1 in
start)
/bin/su $TOMEE_USER $CATALINA_HOME/bin/startup.sh
;; 
stop)   
/bin/su $TOMEE_USER $CATALINA_HOME/bin/shutdown.sh
;; 
restart)
/bin/su $TOMEE_USER $CATALINA_HOME/bin/shutdown.sh
/bin/su $TOMEE_USER $CATALINA_HOME/bin/startup.sh
;; 
esac    
exit 0